# src/files/base/etc/init.d/alpine-persistence
#!/sbin/openrc-run

description="Mount optional Alpine boot and data partitions"

depend() { after dev mdev hwdrivers ; }

start() {
  get_config() { /usr/libexec/alpine-get-config "$1"; }
  local boot_mounted=false data_mounted=false

  ebegin "Mounting Persistence Partitions"
  mkdir -p /boot /data

  # Attempt to mount boot
  if [ -e "/dev/disk/by-partlabel/alpine-boot" ]; then
    if mount -o defaults,noatime /dev/disk/by-partlabel/alpine-boot /boot; then
      einfo "Mounted alpine-boot on /boot"
      boot_mounted=true
    else
      eerror "Failed to mount alpine-boot"
    fi 
  fi

  # Attempt to mount data
  if [ -e "/dev/disk/by-partlabel/alpine-data" ]; then
    if mount -o defaults,noatime /dev/disk/by-partlabel/alpine-data /data; then
      einfo "Mounted alpine-data on /data"
      data_mounted=true
    else
      eerror "Failed to mount alpine-data"
    fi
  fi

  # Only setup persistence if both partitions are mounted
  if [ "$boot_mounted" = true ] && [ "$data_mounted" = true ]; then
    einfo "Setting up persistent files and directories..."
    
    persistent_files=$(get_config persitence.files)
    persistent_dirs=$(get_config persitence.dirs)
    
    for file in $persistent_files; do
      einfo "Setting up persistent file: $file"
      alpine-persist file "$file"
    done

    for dir in $persistent_dirs; do
      einfo "Setting up persistent directory: $dir"
      alpine-persist dir "$dir"
    done
  else
    ewarn "Skipping persistence: One or more partitions is not mounted."
  fi

  eend 0
}

stop() {
  ebegin "Unmounting Persistence Partitions"
  umount /data 2>/dev/null
  umount /boot 2>/dev/null
  eend 0
}