# src/files/base/etc/init.d/alpine-persistence
#!/sbin/openrc-run

description="Mount optional Alpine boot and data partitions"

depend() {
  need dev
  after mdev
}

start() {
  ebegin "Mounting Persistence Partitions"
  mkdir -p /boot /data

  # Attempt to mount boot
  if [ -e "/dev/disk/by-partlabel/alpine-boot" ]; then
    mount -o defaults,noatime /dev/disk/by-partlabel/alpine-boot /boot
    einfo "Mounted alpine-boot on /boot"
  fi

  # Attempt to mount data
  if [ -e "/dev/disk/by-partlabel/alpine-data" ]; then
    mount -o defaults,noatime /dev/disk/by-partlabel/alpine-data /data
    einfo "Mounted alpine-data on /data"
  fi
  
  eend 0
}

stop() {
  ebegin "Unmounting Persistence Partitions"
  umount /data 2>/dev/null
  umount /boot 2>/dev/null
  eend 0
}