# src/files/installer/etc/init.d/installer
#!/sbin/openrc-run
description="Run Alpine Installer"
depend() {
  need localmount net
  after dhcpcd
}

start() {
  # Only run installer if alpine.slot=installer
  slot=$(/usr/bin/alpine-get-cmdline alpine.slot)
  if [ "$slot" = "installer" ]; then
    ebegin "Waiting for DNS and Internet connectivity"    
    local count=0
    while ! wget -q --spider -T 2 https://github.com 2>/dev/null; do
      if [ $count -ge 300 ]; then
        eend 1 "Network/DNS timed out"
        return 1
      fi
      count=$((count + 1)) && sleep 1
    done
    eend 0

    ebegin "Running Alpine Installer"
    alpine-install
    eend $?
  fi
}     